# home-config-autostart/init.desktop
[Desktop Entry]
Name=init
Comment=init
Exec=/etc/osscripts/init
Terminal=false
Type=Application
NoDisplay=true
OnlyShowIn=OPENBOX;GNOME;LXDE
