# home-config-autostart/tmpdir.desktop
[Desktop Entry]
Name=tmpdir
Comment=tmpdir
Exec=mkdir /tmp/Downloads
Terminal=false
Type=Application
NoDisplay=true
OnlyShowIn=OPENBOX;GNOME;LXDE
